# deploy/systemd/cccc-test-memory-palace.service
[Unit]
Description=CCCC Test Memory Palace Backend
After=network.target

[Service]
Type=simple
User=cccc
Group=cccc
ExecStart=/bin/bash -lc 'set -a; source /etc/cccc-test.env; cd "$PROJECT_ROOT/Memory-Palace-main/backend"; exec "$VENV_DIR/bin/python" main.py'
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
