# deploy/systemd/cccc-test-daemon.service
[Unit]
Description=CCCC Test Daemon
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash -lc 'set -a; source /etc/cccc-test.env; cd "$PROJECT_ROOT"; exec "$VENV_DIR/bin/python" scripts/run_cccc_daemon.py'
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
